# services/n20-home-dashboard/service/log/run
#!/data/data/com.termux/files/usr/bin/sh
export PREFIX=/data/data/com.termux/files/usr
logdir="$PREFIX/var/log/sv/n20-home-dashboard"
mkdir -p "$logdir"
exec "$PREFIX/bin/svlogd" -tt "$logdir"
