# templates/services/filebrowser-quantum/run
#!/data/data/com.termux/files/usr/bin/sh
set -u

export HOME="${HOME:-/data/data/com.termux/files/home}"
export PREFIX="${PREFIX:-/data/data/com.termux/files/usr}"
export PATH="$HOME/.local/bin:$HOME/bin:$PREFIX/bin:$PATH"

SERVICE="filebrowser-quantum"
MOUNT="${MOUNT:-/mnt/media_rw/<usb-volume-id>}"
CONFIG="${CONFIG:-$HOME/.config/filebrowser-quantum/config.yaml}"
BIN="${BIN:-$HOME/.local/bin/filebrowser-quantum}"
PASSWORD_FILE="${PASSWORD_FILE:-$HOME/.config/filebrowser-quantum/admin-password.txt}"

fail_retry() {
  printf '%s %s: preflight_failed reason=%s\n' "$(date +%Y-%m-%dT%H:%M:%S%z)" "$SERVICE" "$1"
  sleep 30
  exit 111
}

[ -x "$BIN" ] || fail_retry "binary_missing"
[ -r "$CONFIG" ] || fail_retry "config_missing"
[ -r "$PASSWORD_FILE" ] || fail_retry "password_file_missing"
[ -d "$MOUNT" ] || fail_retry "mount_missing"
[ -r "$MOUNT" ] || fail_retry "mount_not_readable"
ls "$MOUNT" >/dev/null 2>&1 || fail_retry "mount_not_listable"

mkdir -p "$HOME/.local/share/filebrowser-quantum" \
  "$HOME/.cache/filebrowser-quantum" \
  "$HOME/.local/state/filebrowser-quantum"

FILEBROWSER_ADMIN_PASSWORD=$(tr -d '\r\n' < "$PASSWORD_FILE")
export FILEBROWSER_ADMIN_PASSWORD

exec "$BIN" -c "$CONFIG" 2>&1
